# supervisord.conf
[supervisord]
nodaemon=true

[program:dockerd]
user=root
command=dockerd-entrypoint.sh
autorestart=false
startsecs=0